# docker/cron/crontab
* * * * * /usr/local/bin/php /var/www/bin/console app:process-orders >> /var/log/cron.log 2>&1
